# server-config/config/sample-app.service
[Unit]
Description=Sample NodeJS Application

[Service]
ExecStart=
ExecStop=

[Install]
